# ScheduleTimelapse.service
[Unit]
Description=Schedule Timelapse service

[Service]
ExecStart=/home/kolbe/hamcam/schedule_timelapse.bash
Type=oneshot
